# dock_poc/docker/pocsp/conf/supervisor/supervisord_stormsv.conf
[supervisord]
pidfile=/var/run/supervisord.pid
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


[program:stormsv]
command=storm supervisor
redirect_stderr=true
stdout_logfile=/var/log/stormsv.out.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
priority=1

[program:stormlogviewer]
command=storm logviewer
redirect_stderr=true
stdout_logfile=/var/log/stormlogviewer.out.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true




