# dock_poc/docker/pocsp/conf/supervisor/supervisord_yarnrm.conf
[supervisord]
pidfile=/var/run/supervisord.pid
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


[program:yarnrm]
command=watchpid.py /tmp/yarn--resourcemanager.pid %(ENV_HADOOP_YARN_HOME)s/sbin/yarn-daemon.sh start resourcemanager
redirect_stderr=true
stdout_logfile=/var/log/yarnrm.out.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
priority=1


[program:yarnhs]
command=watchpid.py /tmp/mapred--historyserver.pid %(ENV_HADOOP_YARN_HOME)s/sbin/mr-jobhistory-daemon.sh start historyserver
redirect_stderr=true
stdout_logfile=/var/log/yarnhs.out.log
stdout_logfile_maxbytes=10MB
priority=2
autostart=true
autorestart=true

[program:deployjobs]
command=/usr/local/bin/runsamzajobs.sh
redirect_stderr=true
stdout_logfile=/var/log/runsamzajobs.out.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=false


