# dock_poc/docker/pocsp/conf/supervisor/supervisord_yarnnm.conf
[supervisord]
pidfile=/var/run/supervisord.pid
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:yarnnm]
command=watchpid.py /tmp/yarn--nodemanager.pid %(ENV_HADOOP_YARN_HOME)s/sbin/yarn-daemon.sh start nodemanager
redirect_stderr=true
stdout_logfile=/var/log/yarnnm.out.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true



